# systemd/zohar-api.service
[Unit]
Description=Zohar API Dashboard (FastAPI)
After=network.target zohar-llama.service

[Service]
Type=simple
User=gorops
WorkingDirectory=/home/gorops/proyectos antigravity/zohar-agent/api
ExecStart="/home/gorops/proyectos antigravity/zohar-agent/zohar_venv/bin/python" -m uvicorn main:app --host 0.0.0.0 --port 8081
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
