# systemd/zohar-backup.timer
[Unit]
Description=Run Zohar Backup every day at midnight

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
