# systemd/zohar-api.service
[Unit]
Description=Zohar API Dashboard (FastAPI)
After=network.target zohar-llama.service

[Service]
Type=simple
User=gorops
WorkingDirectory=/home/gorops/zohar-agent/api
ExecStart=/home/gorops/zohar-agent/zohar_venv/bin/python zohar_api.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
